# templates/default/service_master-api-containerized.service.erb
[Unit]
Description=Atomic OpenShift Master API
Documentation=https://github.com/openshift/origin
After=etcd_container.service
Wants=etcd_container.service
Before=<%= node['cookbook-openshift3']['openshift_service_type'] %>-node.service
After=docker.service
PartOf=docker.service
Requires=docker.service

[Service]
EnvironmentFile=/etc/sysconfig/<%= node['cookbook-openshift3']['openshift_service_type'] %>-master-api
Environment=GOTRACEBACK=crash
ExecStartPre=-/usr/bin/docker rm -f <%= node['cookbook-openshift3']['openshift_service_type'] %>-master-api
ExecStart=/usr/bin/docker run --rm --privileged --net=host --name <%= node['cookbook-openshift3']['openshift_service_type'] %>-master-api --env-file=/etc/sysconfig/<%= node['cookbook-openshift3']['openshift_service_type'] %>-master-api -v <%= node['cookbook-openshift3']['openshift_data_dir'] %>:<%= node['cookbook-openshift3']['openshift_data_dir'] %> -v /var/run/docker.sock:/var/run/docker.sock -v <%= node['cookbook-openshift3']['openshift_common_base_dir'] %>:<%= node['cookbook-openshift3']['openshift_common_base_dir'] %> <%= node['cookbook-openshift3']['openshift_docker_master_image'] %>:${IMAGE_VERSION} start master api --config=${CONFIG_FILE} $OPTIONS
ExecStartPost=/usr/bin/sleep 10
ExecStop=/usr/bin/docker stop <%= node['cookbook-openshift3']['openshift_service_type'] %>-master-api
LimitNOFILE=131072
LimitCORE=infinity
WorkingDirectory=<%= node['cookbook-openshift3']['openshift_data_dir'] %>
SyslogIdentifier=<%= node['cookbook-openshift3']['openshift_service_type'] %>-master-api
Restart=always
RestartSec=5s

[Install]
WantedBy=docker.service
WantedBy=<%= node['cookbook-openshift3']['openshift_service_type'] %>-node.service
